# puppet/modules/revealjs/templates/revealjs.conf.erb
#/etc/init/revealjs.conf
description "A reveal.js runtime"

env GRUNT_EXEC=/usr/local/bin/grunt
env APP_DIR='<%= @src_dir %>'
env NODE_PATH=$NODE_PATH:<%= @src_dir %>/node_modules
env PID_FILE=/var/run/revealjs.pid
env LOG_FILE=/var/log/revealjs.log

start on (local-filesystems and net-device-up IFACE=eth0) or stopped
stop on shutdown

pre-start script
	cd $APP_DIR
    $GRUNT_EXEC serve 1>$LOG_FILE 2>&1 &
	echo $! > $PID_FILE;
end script

post-stop script
	if [ -f $PID_FILE ]
	then
        PID=`cat $PID_FILE`
        if [ `ps aux | gawk '{print $2}' | grep ^$PID$` ]
        then
            echo "Stopping Reveal.js Node server..."
            echo "Killing $PID"
            kill $PID;
            echo "Reveal.js Node server stopped"
        else
            "Reveal.js process id ${PID} not running, removing PID file."
        fi
    fi
    rm -f $PID_FILE;
    sleep 10
end script